# ops/systemd/galaxia-coder.service
[Unit]
Description=Galaxia — coder agent (open daily improvement PRs from the veille)
Documentation=https://github.com/Jeffchoux/galaxia_os/tree/main/agents/coder
Wants=network-online.target
After=network-online.target galaxia-veille.service

[Service]
Type=oneshot
User=galaxia
Group=galaxia
WorkingDirectory=/home/galaxia/galaxia-project/agents/coder
# Le wizard écrit ANTHROPIC_API_KEY ici. Le préfixe `-` rend la source
# optionnelle : pas d'échec si le fichier n'existe pas encore (l'agent
# fail-closera proprement avec "ANTHROPIC_API_KEY missing").
EnvironmentFile=-/opt/galaxia/config/.env
# GH_TOKEN doit être présent dans ce fichier pour que `gh pr create` réussisse
# sous systemd (ProtectSystem=strict rend ~/.config/gh illisible par le service).
# Ajouter dans /opt/galaxia/config/.env :
#   GH_TOKEN=ghp_xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx
# Scope minimum : `public_repo` (le repo galaxia_os est public).
ExecStart=/usr/bin/node /home/galaxia/galaxia-project/agents/coder/index.mjs
# Cap à 60 min : un run normal tient en 5–15 min. Au-delà = boucle, on coupe.
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal

# Garde-fous
NoNewPrivileges=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# ProtectSystem=strict monte aussi /tmp en read-only par défaut. PrivateTmp=true
# donne au service son propre namespace /tmp writable et isolé des autres services.
PrivateTmp=true
# ProtectSystem=strict monte TOUT le reste (dont /home) en lecture seule. L'agent
# clone et édite dans son /tmp privé, donc il n'a pas besoin d'écrire sous /home…
# SAUF pour consommer les propositions curées : le pont déplace
# galaxia-updates/pending/*.md vers applied/. On ré-ouvre en écriture cet arbre
# précis, et lui seul. Le reste de /home (clés SSH, config gh) reste lisible.
ReadWritePaths=/home/galaxia/.claude/galaxia/galaxia-updates
